# sleep.target
[Unit]
Description=User-level target triggered when the system is about to sleep.
Requires=lock.target
After=lock.target
